# cron/anacrontab
# /etc/anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/sh
HOME=/root
LOGNAME=root

# These replace cron's entries
1	5	cron.daily	run-parts --report /etc/cron.daily
7	10	cron.weekly	run-parts --report /etc/cron.weekly
@monthly	15	cron.monthly	run-parts --report /etc/cron.monthly

1	15	backup_books	/bin/bash /home/avishek/Code/Shell/backup_books.sh
1	25	backup_zotero	/bin/bash /home/avishek/Code/Shell/backup_zotero.sh
1	35	backup_calibre  /bin/bash /home/avishek/Code/Shell/backup_calibre.sh
1	45	backup_tex  /bin/bash /home/avishek/Code/Shell/backup_tex.sh
1	50	backup_web  /bin/bash /home/avishek/Code/Shell/backup_web.sh
1	55	backup_chatgpt  /bin/bash /home/avishek/Code/Shell/backup_chatgpt.sh